# .config/systemd/user/20rule.timer
[Unit]
Description=Run 20-20-20 rule every 20 minutes

[Timer]
OnCalendar=*:0/20
#Persistent=true
AccuracySec=1s
RandomizedDelaySec=0

[Install]
WantedBy=timers.target
