# .config/systemd/user/20rule.service
[Unit]
Description=20-20-20 Reminder

[Service]
Type=exec
ExecStart=/home/odin/Scripts/20rule.sh
